# docker/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=0.0.0.0:9001
username=user
password=pass

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


[include]
files = /etc/supervisor/conf.d/*.conf

[program:sitebase_queue]
command=/var/www/html/bin/console queue:process
autostart=true
autorestart=true
stderr_logfile=/var/log/sitebase_queue.err.log
stdout_logfile=/var/log/sitebase_queue.out.log

[program:sitebase_cron]
command=/bin/bash -c "/var/www/html/bin/console cron:run && sleep 60"
autostart=true
autorestart=true
stderr_logfile=/var/log/sitebase_cron.err.log
stdout_logfile=/var/log/sitebase_cron.out.log